# files/vm-resources/ABANDONED/ubu-droid/services/system/waydroid-ping.service
[Unit]
Description=Waydroid Ping Service

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/waydroid container unfreeze