# files/vm-resources/waydroid/services/user/phoc.service
[Unit]
Description=Phoc

[Service]
Type=simple
ExecStart=/usr/bin/phoc

[Install]
WantedBy=default.target