# scripts/vm/10-vfio/custom-vfio-hook-files/vfio-load.service
[Unit]
Description=Insert vfio-pci driver

[Service]
Type=oneshot
ExecStart=modprobe -i vfio-pci

[Install]
WantedBy=multi-user.target