# files/vm-resources/ubu-droid/services/user/waydroid-session.service
[Unit]
Description=Waydroid Service

[Service]
Type=simple
ExecStart=/home/android/bin/waydroid-start

[Install]
WantedBy=default.target